# service/nginx/run
#!/bin/bash

# Requires 'daemon off;' in nginx.conf

echo Starting NGINX
exec 2>&1
exec ../../../../apps/sbin/nginx

