# service/gunicorn-django/run
#!/bin/bash

echo Starting gunicorn_django
exec 2>&1
source ../tools/python_virtualenv/bin/activate
cd ../application/django
exec gunicorn_django


