# service/nodejs/run
#!/bin/sh

echo Starting NodeJS
exec 2>&1
exec envdir ./env node ../../bundle/main.js

